# puppet/modules/redis/templates/sv/log/run
#!/bin/bash

set -e
exec 2>&1
exec svlogd -tt -b 10000 .